# node-exporter.service
[Unit]
Description=Node Exporter Prometheus

[Service]
Type=oneshot
ExecStart=/usr/bin/node_exporter

[Install]
WantedBy=multi-user.target